# ansible/roles/log_stats/templates/log-stats.timer.j2
[Unit]
Description=Run log-stats.service every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=log-stats.service

[Install]
WantedBy=timers.target
